# .config/systemd/user/copyq.service
[Unit]
Description=CopyQ Clipboard Manager

[Service]
ExecStart=/usr/bin/copyq

[Install]
WantedBy=default.target
